# templates/systemd.service.epp
# Managed by Puppet, DO NOT EDIT

[Unit]
Description=Renew dehydrated certificates about to expire

[Service]
Type=oneshot
ExecStart=<%= $dehydrated::bin %> --accept-terms --cron --keep-going
User=<%= $dehydrated::user %>
Group=<%= $dehydrated::user %>
